# deployment/97layeros-nightguard.service
[Unit]
Description=97LAYER Night Guard (정찰기)
Documentation=https://github.com/97layer/97layerOS
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/97layerOS
ExecStart=/usr/bin/python3 /home/ubuntu/97layerOS/execution/system/nightguard_daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/97layer-nightguard.log
StandardError=append:/var/log/97layer-nightguard.log

# 환경변수
Environment="ENVIRONMENT=GCP_VM"
Environment="PROCESSING_MODE=sequential"
Environment="PYTHONUNBUFFERED=1"

# 리소스 제한 (RAM 700MB)
MemoryMax=700M
MemoryHigh=600M

[Install]
WantedBy=multi-user.target
