# knowledge/docs/deployment/97layer-nightguard.service
[Unit]
Description=97layerOS Nightguard V2 — Self-Diagnostic Autonomous Daemon
Documentation=https://github.com/97layer/97layerOS
After=network.target 97layer-telegram.service
Wants=97layer-telegram.service

[Service]
Type=simple
User=USERNAME_PLACEHOLDER
WorkingDirectory=/home/USERNAME_PLACEHOLDER/97layerOS
Environment="PYTHONPATH=/home/USERNAME_PLACEHOLDER/97layerOS"
Environment="PATH=/home/USERNAME_PLACEHOLDER/97layerOS/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/USERNAME_PLACEHOLDER/97layerOS/.env
ExecStart=/home/USERNAME_PLACEHOLDER/97layerOS/.venv/bin/python3 core/system/nightguard_v2.py
Restart=always
RestartSec=30

# GCP e2-micro 메모리 제약 (1GB RAM)
MemoryMax=150M
MemoryHigh=100M

# CPU 스로틀링 (백그라운드 데몬)
CPUQuota=20%

# 로깅
StandardOutput=append:/home/USERNAME_PLACEHOLDER/97layerOS/.infra/logs/nightguard.log
StandardError=append:/home/USERNAME_PLACEHOLDER/97layerOS/.infra/logs/nightguard.error.log

[Install]
WantedBy=multi-user.target
