# knowledge/docs/deployment/97layer-code-agent.service
[Unit]
Description=97layer Code Agent — 텔레그램 코드 수정 자동화
After=network.target 97layer-telegram.service
Wants=97layer-telegram.service

[Service]
ExecStart=python3 /home/skyto5339_gmail_com/97layerOS/core/agents/code_agent.py --schedule
WorkingDirectory=/home/skyto5339_gmail_com/97layerOS
EnvironmentFile=/home/skyto5339_gmail_com/97layerOS/.env
Restart=on-failure
RestartSec=15
User=skyto5339_gmail_com
MemoryMax=200M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
