# knowledge/docs/deployment/97layer-ecosystem.service
[Unit]
Description=97layerOS Organic Ecosystem — THE CYCLE Full Stack
Documentation=https://github.com/97layer/97layerOS
After=network.target

[Service]
Type=simple
User=USERNAME_PLACEHOLDER
WorkingDirectory=/home/USERNAME_PLACEHOLDER/97layerOS
Environment="PYTHONPATH=/home/USERNAME_PLACEHOLDER/97layerOS"
Environment="PATH=/home/USERNAME_PLACEHOLDER/97layerOS/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/USERNAME_PLACEHOLDER/97layerOS/.env
ExecStart=/bin/bash /home/USERNAME_PLACEHOLDER/97layerOS/start_ecosystem.sh
Restart=on-failure
RestartSec=15

# GCP e2-micro 메모리 제약
# ecosystem = heartbeat + signal_router + scheduler + SA + AD + CE (6개 프로세스)
# 총 메모리 한도: 600M (각 100M 기준)
MemoryMax=600M
MemoryHigh=500M

# 로깅
StandardOutput=append:/home/USERNAME_PLACEHOLDER/97layerOS/.infra/logs/ecosystem.log
StandardError=append:/home/USERNAME_PLACEHOLDER/97layerOS/.infra/logs/ecosystem.error.log

[Install]
WantedBy=multi-user.target
